# supervisor/gitbucket.conf

[program:gitbucket]
directory=/
command=/usr/bin/java -jar /gitbucket.war
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_error.log
